# install/linux/check.cron
#!/bin/bash

PID=`pgrep -u root usbmond`

if [ -z "${PID}" ]
then
 /opt/usbmond/bin/usbmond >> /dev/null 
fi

